# packaging/remotelauncher.service
[Unit]
Description=RemoteLauncher — launch desktop apps from Android
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Без этих переменных запущенные программы не смогут открыть окно
# (Wayland сокет в %t/wayland-1, DBus в %t/bus, X11 для XWayland-приложений).
# WAYLAND_DISPLAY=wayland-1 — стандартное имя сокета большинства композиторов.
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_SESSION_TYPE=wayland
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
ExecStart=%h/.local/bin/remotelauncher
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
